# test/Demo/src/docker-files/supervisor/conf.d/consul.conf
[program:consul]
directory=/consul
# docker run -i --network lab_consul01 
# -v D:\temp\consuldemo\singleconsul\data\:/mnt/data/consul 
# -v D:\temp\consuldemo\singleconsul\consul.config:/etc/config/consul.json
# consul agent -config-file /etc/config/consul.json
command=/consul/consul agent -config-file=/etc/config/consul.json
autostart=true
