# test/Demo/src/docker-files/supervisord/supervisord.d/app.conf
[program:app]
directory=/app
command=dotnet app.dll
autostart=true
